# php_nginx/container-files/etc/supervisor.d/beanstalkd.conf
[program:beanstalkd]
command = beanstalkd -b /var/lib/beanstalkd
autorestart = true
stdout_logfile = /data/logs/beanstalkd.log
stderr_logfile = /data/logs/beanstalkd.log
